# target/availability_report.service
[Unit]
Description=Availability Report

[Service]
ExecStart=/usr/local/sbin/availability_report

[Install]
WantedBy=default.target
